# ops/observability/systemd/catan-fleet-exporter.service
[Unit]
Description=Catan-Zero generation Prometheus exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/catan-zero-v1/.venv/bin/python /home/ubuntu/catan-zero-v1/tools/fleet/fleet_metrics_exporter.py --listen 127.0.0.1 --port 9500 --run-root /home/ubuntu/gen_out
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
